# containers/systemd/qem-dashboard.container
[Container]
Pod=qem.pod
Image=qem-dashboard.build
